# scripts/supervisor.d/minio.ini
[program:minio]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
user=minio
command=minio server /var/lib/minio/data --console-address ":9001"