# scripts/supervisor.d/postgres.ini
[program:postgres]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
user=postgres
command=logger postgres -D /var/lib/postgresql/data