# scripts/supervisor.d/redis.ini
[program:redis]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
user=redis
command=logger redis-server /etc/redis/redis.conf