# scripts/template.ini
[program:foo]
autostart=true
command=echo "Hello, World!"
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout



